# server/config/mezzofy-beat.service
[Unit]
Description=Mezzofy AI Celery Beat Scheduler
After=network.target redis.service mezzofy-celery.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mz-ai-assistant/server
EnvironmentFile=/home/ubuntu/mz-ai-assistant/server/config/.env
Environment=PATH=/home/ubuntu/mz-ai-assistant/server/venv/bin
ExecStart=/home/ubuntu/mz-ai-assistant/server/venv/bin/celery -A app.tasks.celery_app beat --loglevel=info --scheduler app.tasks.beat_schedule:DatabaseScheduler
Restart=always
RestartSec=10
SyslogIdentifier=mezzofy-beat

[Install]
WantedBy=multi-user.target
